# docker/supervisor/conf.d/nginx.conf
[program:nginx]
name = nginx
command = nginx -g 'daemon off;'
autostart = true
autorestart = false
startsecs = 1
startretries = 3
exitcodes = 0,2
stopsignal = QUIT
stopwaitsecs = 10